# GSmail.conf
[program:GSmail]
command=python3 /home/curator4/GSmail/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/GSmail/test.err.log
stdout_logfile=/home/curator4/GSmail/test.out.log